# INSTALL/root/etc/profile.d/crypt-bufferdisk.sh
[[ $- == *i* ]] || return 0

/opt/AmutaQ!/crypt/cryptdisk_open_bufferdisk
