# INSTALL/root/etc/profile.d/crypt-bufferdisk.sh
[[ $- == *i* ]] || return 0

/opt/AmutaQ!/tools/mount_bufferdisk
